# Homebridge/homebridge.service
[Unit]
Description=Homebridge

[Service]
Type=simple
User=root
ExecStart=/home/*/.homebridge/start.sh

[Install]
WantedBy=multi-user.target
